# jellyfin-broadcaster.service
[Unit]
Description=Jellyfin Discovery Server
After=network.target

[Service]
Type=simple
# Uncomment and set User to run as a specific user (requires network permissions)
# User=your-username
ExecStart=/usr/local/bin/jellyfin-broadcaster -server-url "http://192.168.1.100:8096" -server-id "your-server-id-here" -server-name "Jellyfin Server"
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jellyfin-broadcaster

[Install]
WantedBy=multi-user.target

